# gorilla.cron
0 1,13 * * * root /bin/sh /opt/psql-email-alert/main.sh
